# deb/DEBIAN/prerm
#!/bin/bash
if [ -f /etc/systemd/system/discord-updater.service ]; then
    systemctl stop discord-updater.service
    systemctl disable discord-updater.service
fi